# deploy/systemd/privacynode-issuer.service
[Unit]
Description=Privacy Node Issuer Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=privacynode
Group=privacynode
WorkingDirectory=/var/lib/privacynode
EnvironmentFile=-/etc/privacynode/common.env
EnvironmentFile=-/etc/privacynode/issuer.env
ExecStart=/usr/local/bin/node --issuer
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/privacynode
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
